# bot/trade_trigger/ops/qa-backup.service
[Unit]
Description=QuantumAlpha SQLite backup
After=network-online.target

[Service]
Type=oneshot
User=qa
Group=qa
ExecStart=/home/qa/quantumalpha/ops/backup.sh
StandardOutput=append:/home/qa/backups/backup.log
StandardError=append:/home/qa/backups/backup.log
